# root/etc/s6-overlay/s6-rc.d/init-permissions-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# permissions
find /config/log ! -path '/config/log/logrotate.status' -exec chmod +r {} \+

lsiown -R abc:abc \
    /config
