# root/etc/s6-overlay/s6-rc.d/init-swag-folders/up
/etc/s6-overlay/s6-rc.d/init-swag-folders/run
